# supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700
username=admin
password=strongpassword

[supervisord]
logfile=/var/log/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=true
user=root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

# Nginx Web Server
[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/nginx.log
stderr_logfile=/var/log/nginx.error.log
user=root

# FastAPI Main Application
[program:bridge-api]
command=uvicorn backend.main:app --host 0.0.0.0 --port 8000 --workers 2
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/bridge-api.log
stderr_logfile=/var/log/bridge-api.error.log
user=app
environment=PYTHONPATH="/app"

# Cascoin Watcher Service
[program:cascoin-watcher]
command=python3 /app/watchers/cascoin_watcher.py
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/cascoin-watcher.log
stderr_logfile=/var/log/cascoin-watcher.error.log
user=app
environment=PYTHONPATH="/app" 